# stack/config/sonarqube/system/profile.d/postgres.sh
export POSTGRES_HOME=/usr/share/postgresql

export PATH=${POSTGRES_HOME}/bin:${PATH}
