# stack/config/sonarqube/system/profile.d/shared-lib.sh
export LD_LIBRARY_PATH=/usr/share/postgesql/lib
