# stack/config/sonarqube/system/profile.d/sonarqube.sh
export SONARQUBE_HOME=/usr/share/sonarqube
export SONARQUBE_SERVICE=${SONARQUBE_HOME}/sonarqube-7.9/bin/linux-x86-64

export PATH=${PATH}:${SONARQUBE_SERVICE}